# exp_images/jenkins/supervisor_jenkins.conf
[program:jenkins]
command=service jenkins start
startsec=0
autorestart=false
